# root/etc/supervisor2/conf.d/crond.conf
[program:crond]
command=/usr/sbin/cron -f
